# feed/configuration-update-mt7620/uci-defaults/01_initial_boot_start
#!/bin/ash

. /bin/ledcontrol.ash
led_quick_flashing
