# feed/configuration-update-ib2/uci-defaults/01_initial_boot_start
#!/bin/ash

. /bin/ledcontrol.ash
led_info_heartbeat
