# feed/configuration-update-ib2/uci-defaults/99_initial_boot_end
#!/bin/ash

sync
. /bin/ledcontrol.ash
led_info_off
