# feed/configuration-update-mt7621/uci-defaults/01_initial_boot_start
#!/bin/ash

. /bin/ledcontrol.ash
led_boot_start
