# feed/configuration-update-mt7621/uci-defaults/99_initial_boot_end
#!/bin/ash

sync
. /bin/ledcontrol.ash
led_boot_end
